# ansible/templates/docker-compose.service
[Unit]
Description={{ application.name }} with Docker
Requires=docker.service

[Service]
User={{ application.name }}
Group={{ application.name }}
Restart=always
ExecStart=/usr/bin/docker compose -f /opt/docker-compose/{{ application.name }}/docker-compose.yml up
ExecStop=/usr/bin/docker compose -f /opt/docker-compose/{{ application.name }}/docker-compose.yml down

[Install]
WantedBy=default.target
